# com.mediadash.cache-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mediadash.cache-sync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mattfasullo/Projects/MediaDash/sync_shared_cache.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>1800</integer>
    <!-- Run every 30 minutes (1800 seconds) -->
    
    <key>RunAtLoad</key>
    <true/>
    <!-- Run immediately when loaded -->
    
    <key>StandardOutPath</key>
    <string>/tmp/mediadash-cache-sync.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/mediadash-cache-sync-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>KeepAlive</key>
    <false/>
    <!-- Don't restart if it crashes - just wait for next interval -->
    
    <key>ThrottleInterval</key>
    <integer>300</integer>
    <!-- Don't restart more than once every 5 minutes if it crashes -->
</dict>
</plist>

